# scripts/claude-trader.service
[Unit]
Description=Claude Code Trader (tmux persistent session)
After=network-online.target redis-server.service
Wants=network-online.target

[Service]
Type=forking
User=root
WorkingDirectory=/root/Projects/ai-trading-bot
Environment=HOME=/root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/root/.npm-global/bin
ExecStartPre=-/usr/bin/tmux kill-session -t claude-trader
ExecStart=/usr/bin/tmux new-session -d -s claude-trader -c /root/Projects/ai-trading-bot 'claude --continue'
ExecStop=/usr/bin/tmux kill-session -t claude-trader
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
